# deploy/quadlet/dashboard_pi.container
[Unit]
Description=Dashboard PI (Flask) - prod
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=dashboard_pi
# ⚠️  Reemplaza <usuario> y <DIGEST> con el valor del job summary de GitHub Actions
Image=ghcr.io/proyectopi/dashboard_pi@sha256:5f13ad51fd959938d7ebc6662c194c41995f4a652f1bb79e4aba3055cf042152
Network=internal.network
Volume=dashboard_pi_data:/data:z
EnvironmentFile=/root/dashboard_pi.env
# AuthFile=/etc/containers/auth.json   # SOLO si el paquete GHCR es privado (ver Fase 4.2)
PodmanArgs=--memory=384m --cpus=0.6

[Service]
Restart=always

[Install]
WantedBy=multi-user.target default.target
